# tf/squirrelcart_backup.service
[Unit]
Description=Backup Squirrelcart every 30 mins

[Service]
Type=oneshot
ExecStart=/bin/bash /home/centos/squirrel/backup_squirrel.sh
