# tf/squirrelcart_cert.timer
[Unit]
Description=Check for letsencrypt cert update every Friday

[Timer]
OnCalendar=Sat *-*-* 18:00:00

[Install]
WantedBy=timers.target
